# overlay/lower/etc/network/if-up.d/odhcp6c
#!/bin/sh

PIDFILE="/var/run/odhcp6c_$IFACE.pid"
TIMEOUT=10

. /usr/share/common

if [ "$(sysctl -n net.ipv6.conf.all.disable_ipv6)" -eq 1 ]; then
	echo "$IFACE: IPv6 disabled globally"
	exit 0
fi

if [ "true" = "$IF_DHCP_V6_ENABLED" ]; then
	echo "$IFACE: DHCPv6 enabled"
else
	echo "$IFACE: DHCPv6 disabled"
	exit 0
fi

do_start() {
	TMP_LOG=/tmp/odhcp6c.log
	/sbin/odhcp6c -p "$PIDFILE" -t120 -S15 "$IFACE" > "$TMP_LOG" 2>&1 &

	while [ "$TIMEOUT" -gt 0 ]; do
		if grep -q "Adding address" "$TMP_LOG"; then
			echo "$IFACE: DHCPv6 configured"
			rm -f "$TMP_LOG"
			return 0
		fi
		sleep 1
		TIMEOUT=$((TIMEOUT - 1))
	done

	echo "$IFACE: DHCPv6 timeout reached"
	rm -f "$TMP_LOG"
	return 1
}

do_stop() {
	if [ -f "$PIDFILE" ]; then
		kill "$(cat "$PIDFILE")"
		rm -f "$PIDFILE"
	fi
}

case "$MODE" in
	start)
		case "$PHASE" in
			post-up)
				echo "$IFACE: starting odhcp6c"
				do_start
				;;
		esac
		;;

	stop)
		case "$PHASE" in
			pre-down)
				echo "$IFACE: stopping odhcp6c"
				do_stop
				;;
		esac
		;;

	*)
		echo "Usage: $0 {start|stop}"
		exit 1
		;;
esac

exit 0
