# overlay/lower/etc/network/if-pre-up.d/disable_wlan
#!/bin/sh

. /etc/init.d/rc.common

if [ "$IFACE" = "wlan0" ]; then
	if [ -z "$(get wlanssid)" ] || [ -z "$(get wlanpass)" ]; then
		info "wlan0 disabled, no SSID/pass"
		exit 1
	elif [ "$(get disable_wlan)" = "true" ]; then
		info "WLAN manually disabled."
		exit 1
	fi
fi
