# overlay/lower/etc/network/if-pre-up.d/disable_wlan
#!/bin/sh

. /etc/init.d/rc.common

if [ "$IFACE" = "wlan0" ] && { [ -z "$(get wlanssid)" ] || [ -z "$(get wlanpass)" ]; }; then
	info "wlan0 disabled, no ssid/pass"
	exit 1
fi
