# overlay/lower/etc/network/if-pre-up.d/disable_wlan
#!/bin/sh

# only work on a WLAN interface
if [ "wlan0" != "$IFACE" ]; then
	exit 0
fi

. /usr/share/common

wlan_pass=$(jct /etc/thingino.json get wlan.pass 2>/dev/null || fw_printenv -n wlan_pass)
wlan_ssid=$(jct /etc/thingino.json get wlan.ssid 2>/dev/null || fw_printenv -n wlan_ssid)
wlan_enabled=$(jct /etc/thingino.json get wlan.enabled 2>/dev/null || fw_printenv -n wlan_enabled)

# exit 1 to stop the interface from being brought up, and dhcp being attempted
if [ -z "$wlan_ssid" ]; then
	echo_error "WLAN disabled, no SSID"
	exit 1
elif [ -z "$wlan_pass" ]; then
	echo_error "WLAN disabled, no password"
	exit 1
elif [ "$wlan_enabled" = "false" ]; then
	echo_error "WLAN manually disabled."
	exit 1
fi

exit 0
