# overlay/lower/etc/network/if-pre-up.d/set-mac-address
#!/bin/sh

. /etc/init.d/rc.common

skip="Skip MAC address setup"

[ "eth0" = "$IFACE" ] || \
	quit "$skip $IFACE is not eth0"

nfsip=$(sed -n '/^[^ ]*:.* \/ nfs.*[ ,]addr=\([0-9.]\+\).*/s//\1/p' /proc/mounts)

[ -n "$nfsip" ] && ip route get "$nfsip" | grep -q "dev eth0" && \
	quit "$skip, $IFACE used for NFS from $nfsip"

[ -d /sys/class/net/eth0 ] || \
	quit "$skip, eth0 does not exist"

readlink -f /sys/class/net/eth0/device | grep -q "usb" && \
	quit "$skip, eth0 is a USB device"

MAC_ADDR=$(fw_printenv -n ethaddr)

[ -z "$MAC_ADDR" ] && \
	quit "$skip, ethaddr not found in U-Boot environment"

/sbin/ip link set dev eth0 address "$MAC_ADDR" && \
	quit "Set eth0 MAC address from environment: $MAC_ADDR"

die "Failed to set eth0 MAC address from environment: $MAC_ADDR"
