# overlay/lower/etc/crontabs/root
# see https://github.com/themactep/thingino-firmware/wiki/Cron-jobs

# reboot camera nightly @ 3:00
#0 3 * * * reboot -f
# run daynight every minute
*/1 * * * * daynight
