# overlay/etc/rc.local
#!/bin/sh
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.

ircut.sh ~
sleep 0.5
ircut.sh ~
sleep 0.5
daynight.sh day

exit 0
